# keyboard/services/qwamos-pq-keystore.service
[Unit]
Description=QWAMOS Post-Quantum Keystore Service
Documentation=https://github.com/Dezirae-Stark/QWAMOS
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/qwamos/keyboard/crypto
ExecStart=/usr/bin/python3 /opt/qwamos/keyboard/crypto/pq_keystore_service.py --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/local/tmp/qwamos_keystore
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=1024
MemoryMax=128M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
