# security/systemd/qwamos-ml-network-anomaly.service
[Unit]
Description=QWAMOS ML Network Anomaly Detector
Documentation=https://github.com/Dezirae-Stark/QWAMOS
After=network.target qwamos-tor.service qwamos-ai-manager.service
Wants=qwamos-ai-manager.service

[Service]
Type=simple
User=qwamos
Group=qwamos
WorkingDirectory=/opt/qwamos/security/ml

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="TF_CPP_MIN_LOG_LEVEL=2"

# ML Network Anomaly Detector
ExecStart=/usr/bin/python3 /opt/qwamos/security/ml/network_anomaly_detector.py --interface any --threshold 0.15

# Restart on failure
Restart=on-failure
RestartSec=10s

# Resource limits
MemoryLimit=2G
MemoryHigh=1.5G
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/qwamos/security/ml/models /var/log/qwamos
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Network access (needed for packet capture)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_PACKET
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qwamos-ml-network

[Install]
WantedBy=multi-user.target
